# Meta/scripts/com.vault.sprint-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vault.sprint-worker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>HOME_DIR_PLACEHOLDER/VaultSandbox/Meta/scripts/run-sprint-worker.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>HOME_DIR_PLACEHOLDER/VaultSandbox</string>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>StandardOutPath</key>
    <string>/tmp/vault-sprint-worker.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/vault-sprint-worker.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>HOME_DIR_PLACEHOLDER</string>
    </dict>
</dict>
</plist>
